# deploy/baw-telegram.service
[Unit]
Description=BAW Bot — Telegram Messenger Connector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=radxa
Group=radxa
WorkingDirectory=/home/radxa/baw
EnvironmentFile=/home/radxa/.baw/telegram.env
Environment=PATH=/home/radxa/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/radxa/.hermes/hermes-agent/venv/bin/python3 -u /home/radxa/baw/baw-bot --platform telegram --debug
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
